# api/crontab
DJANGO_SETTINGS_MODULE=meus_lares.settings

0 20 * * * cd /app && /usr/local/bin/poetry run python meus_lares/daily_tasks.py
